# ops/systemd/ny-bot.service
[Unit]
Description=Ny Paid Signing Discord Bot
After=network.target

[Service]
User=ny
WorkingDirectory=/opt/ny-paid-signing-fixed
Environment="PATH=/opt/ny-paid-signing-fixed/.venv/bin"
ExecStart=/opt/ny-paid-signing-fixed/.venv/bin/python bot/discord_bot.py
Restart=always

[Install]
WantedBy=multi-user.target
